# 51-jlink_dock.rules
SUBSYSTEM=="tty", ATTRS{idVendor}=="1366", RUN+="/path/to/jlink_dock.sh %E{ACTION} %E{DEVNAME} %E{ID_SERIAL_SHORT}"
